# deploy/tg-bot-daily-report.service.example
[Unit]
Description=Generate Telegram AI daily news report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=tgbot
Group=tgbot
UMask=0077
EnvironmentFile=/etc/tg-bot.env
# DAILY_REPORT_SECTIONS=weather,exchange,market,china,global,ai_tech,proxy,hackernews,github,steam,cold_knowledge
# DAILY_REPORT_STATE_FILE=/var/lib/morning-report/daily_report_state.json
# DAILY_REPORT_STATUS_FILE=/var/lib/morning-report/daily_report_status.json
WorkingDirectory=/opt/tg-bot-search-assistant
ExecStart=/usr/bin/python3 /opt/tg-bot-search-assistant/scripts/build-daily-report.py
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/morning-report
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
