# deploy/tg-bot.service.example
[Unit]
Description=Telegram AI Search Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tgbot
Group=tgbot
UMask=0077
EnvironmentFile=/etc/tg-bot.env
Environment=TG_BOT_DATA_DIR=/var/lib/morning-report
Environment=ASK_API_HOST=127.0.0.1
WorkingDirectory=/opt/tg-bot-search-assistant
ExecStart=/usr/bin/python3 /opt/tg-bot-search-assistant/scripts/tg-bot-new.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/morning-report
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
